# fantastic-probe-monitor.service
[Unit]
Description=ISO MediaInfo Extractor Monitor Service
Documentation=https://github.com/aydomini/fantastic-probe
After=network.target

[Service]
Type=simple
User=root
Group=root

# 服务主程序
ExecStart=/usr/local/bin/fantastic-probe-monitor

# 自动重启策略
Restart=always
RestartSec=10

# 资源限制
Nice=10
IOSchedulingClass=idle

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fantastic-probe

# 安全选项
NoNewPrivileges=false
PrivateTmp=false

[Install]
WantedBy=multi-user.target
