# debian/DEBIAN/postinst
#!/bin/bash
set -e

# postinst script for fantastic-probe

case "$1" in
    configure)
        # 创建配置目录
        mkdir -p /etc/fantastic-probe

        # 如果配置文件不存在，运行配置向导
        if [ ! -f /etc/fantastic-probe/config ]; then
            echo "Fantastic-Probe 已安装！"
            echo ""
            echo "请配置服务："
            echo "  1. 编辑配置文件: /etc/fantastic-probe/config"
            echo "     （已提供模板：/usr/share/doc/fantastic-probe/config.template）"
            echo ""
            echo "  2. 或运行交互式配置："
            echo "     /usr/share/fantastic-probe/configure.sh"
            echo ""
            echo "  3. 启动服务："
            echo "     systemctl enable --now fantastic-probe-monitor"
        else
            echo "配置文件已存在，保留现有配置"
        fi

        # 创建日志文件
        touch /var/log/fantastic_probe.log
        touch /var/log/fantastic_probe_errors.log
        chmod 644 /var/log/fantastic_probe.log
        chmod 644 /var/log/fantastic_probe_errors.log

        # 重新加载 systemd
        systemctl daemon-reload || true
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
