# debian/DEBIAN/prerm
#!/bin/bash
set -e

# prerm script for fantastic-probe

case "$1" in
    remove|upgrade|deconfigure)
        # 停止服务
        if systemctl is-active --quiet fantastic-probe-monitor.service; then
            systemctl stop fantastic-probe-monitor.service || true
        fi

        # 禁用服务
        if systemctl is-enabled --quiet fantastic-probe-monitor.service; then
            systemctl disable fantastic-probe-monitor.service || true
        fi
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
